# .config/systemd/user/ctags.service
[Unit]
Description=Regenerate ctags files

[Service]
ExecStart=/home/marcec/bin/gen_ctags.sh
