# .config/systemd/user/ctags.timer
[Unit]
Description=Regenerate ctags files (timer)

[Timer]
OnCalendar=Sat *-*-* 03:00
Persistent=true

[Install]
WantedBy=default.target
