# .config/systemd/user/getmail@.service
[Unit]
Description=Run getmail with configuration file getmailrc.%I

[Service]
ExecStart=/usr/bin/getmail -n -q --rcfile=getmailrc.%I
