# .config/systemd/user/mpd.service
.include /usr/lib/systemd/system/mpd.service

[Install]
WantedBy=default.target
