# .config/systemd/user/newsbeuter.service
[Unit]
Description=Run newsbeuter -x reload

[Service]
Environment=XDG_DATA_HOME=%h/.config
ExecStart=/usr/bin/newsbeuter -x reload
